# crontab
0 0 * * * root ./manage.py fetch_participant_resources >> /var/log/cron.log
#
